# etc/rc.local
#!/bin/sh -e
#
# rc.local

systemctl enable setup-data.service

exit 0
