# etc/systemd/system/axp20x.service
[Unit]
Description=Configure ax20x PMU
After=mnt-data.mount

[Service]
ExecStart=/etc/bin/config-axp20x.sh

[Install]
WantedBy=mnt-data.mount
