# etc/systemd/system/drop-local-wifi.service
[Unit]
Description=Force reconnect WiFi & ethernet

[Service]
Type=oneshot
ExecStart=/etc/bin/drop-local-wifi.sh
